# packaging/debian/postinst
#!/bin/bash

echo "create symlinks and log folder"
if [ ! -d /var/log/balloon ]; then
  mkdir /var/log/balloon
fi

if [ ! -h /usr/share/balloon/config ]; then
    ln -s /etc/balloon /usr/share/balloon/config
fi

if [ ! -h /usr/share/balloon/log ]; then
    ln -s /var/log/balloon /usr/share/balloon/log
fi

if [ ! -h /usr/bin/ballooncli ]; then
    ln -s /usr/share/balloon/bin/console/ballooncli /usr/bin/ballooncli
fi

if [[ ! -f /etc/systemd/system/balloon-jobs.service && -d /etc/systemd/system ]]; then
    cp /usr/share/balloon/scripts/balloon-jobs.service.systemd /etc/systemd/system/balloon-jobs.service
fi

if [[ ! -f /etc/init/balloon-jobs.conf && -d /etc/init ]]; then
    cp /usr/share/balloon/scripts/balloon-jobs.service.upstart /etc/init/balloon-jobs.conf
fi

if [ ! -d /var/log/balloon ]; then
  mkdir /var/log/balloon
fi

chown -R www-data. /var/log/balloon

if [ ! -d /etc/ssl/balloon ]; then
  mkdir /etc/ssl/balloon
fi

if [[ ! -f /etc/ssl/balloon/key.pem && ! -f /etc/ssl/balloon/chain.pem ]]; then
  openssl genrsa -des3 -passout pass:balloon -out server.pass.key 2048
  openssl rsa -passin pass:balloon -in server.pass.key -out key.pem
  rm server.pass.key
  openssl req -new -key key.pem -out server.csr -subj "/C=CH/L=Zurich/O=Balloon/CN=balloon.local"
  openssl x509 -req -days 365 -in server.csr -signkey key.pem -out chain.pem
  rm server.csr
  mv key.pem /etc/ssl/balloon
  mv chain.pem /etc/ssl/balloon
fi

if [ ! -d /etc/nginx/conf.d/balloon ]; then
  mkdir /etc/nginx/conf.d/balloon
fi

if [ ! -f /etc/nginx/conf.d/balloon.conf ]; then
  cp /usr/share/balloon-web/nginx/nginx.conf /etc/nginx/conf.d/balloon.conf
fi

if [ ! -f /etc/nginx/conf.d/balloon/web.conf ]; then
  cp /usr/share/balloon-web/nginx/nginx-web.conf /etc/nginx/conf.d/balloon/web.conf
fi

echo "execute balloon upgrade"
/usr/bin/ballooncli upgrade start -vvv -i
chown -R www-data. /var/log/balloon

echo "create systemctl service"
systemctl enable balloon-jobs.service 2> /dev/null

echo "restart balloon-jobs service"
if service --status-all | grep -Fq 'balloon-jobs'; then
  service balloon-jobs restart
fi

echo "restart nginx"
service nginx restart
