# packaging/balloon-jobs.service
[Unit]
Description=Ballooon
After=network.target

[Service]
ExecStart=/usr/bin/ballooncli jobs -d
ExecReload=/bin/kill -HUP $MAINPID
User=balloon
KillMode=process
Restart=on-failure
RestartPreventExitStatus=255
Type=simple

[Install]
WantedBy=multi-user.target
Alias=balloon-jobs.service
